# tv-garden-daily.service
[Unit]
Description=IPTV Playlist Daily Update
After=network.target
Wants=network.target

[Service]
Type=oneshot
User=root
Group=root
WorkingDirectory=/home/loganlatham/JellyGarden
ExecStart=/usr/bin/python3 /home/loganlatham/JellyGarden/iptv_extractor.py --output-m3u /var/lib/tv-garden/iptv_playlist.m3u --output-xmltv /var/lib/tv-garden/epg_guide.xmltv
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/tv-garden /var/lib/tv-garden /home/loganlatham/JellyGarden
RestrictNamespaces=true
